# deploy/linux/systemd/hart-agent-daemon.service
[Unit]
Description=HART OS Agent Daemon (Goal Engine)
Documentation=https://github.com/hevolve-ai/hart
After=hart-backend.service
BindsTo=hart-backend.service

[Service]
Type=simple
User=hart
Group=hart
WorkingDirectory=/opt/hart
ExecStart=/opt/hart/venv/bin/python -c "from integrations.agent_engine.agent_daemon import AgentDaemon; d = AgentDaemon(); d.run_forever()"
EnvironmentFile=/etc/hart/hart.env
Restart=on-failure
RestartSec=30
StartLimitBurst=3
StartLimitIntervalSec=60
TimeoutStartSec=30

NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/hart /var/log/hart /opt/hart/agent_data
PrivateTmp=yes

ProtectClock=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
SystemCallFilter=@system-service
MemoryDenyWriteExecute=no
LockPersonality=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes

StandardOutput=journal
StandardError=journal
SyslogIdentifier=hart-agent-daemon

[Install]
WantedBy=hart.target
