# deploy/linux/systemd/hart-llm.service
[Unit]
Description=HART OS Local LLM (llama.cpp)
Documentation=https://github.com/hertz-ai/HARTOS
After=network.target
ConditionPathExists=/opt/hart/models/default.gguf

[Service]
Type=simple
User=hart
Group=hart
ExecStart=/opt/hart/bin/llama-server --model /opt/hart/models/default.gguf --port ${HART_LLM_PORT:-8080} --ctx-size 4096 --threads ${HART_LLM_THREADS:-4}
EnvironmentFile=/etc/hart/hart.env
Restart=on-failure
RestartSec=10
StartLimitBurst=3
StartLimitIntervalSec=60
TimeoutStartSec=60

# LLM needs GPU access
SupplementaryGroups=video render
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/hart/models
PrivateTmp=yes

ProtectClock=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
SystemCallFilter=@system-service
MemoryDenyWriteExecute=no
LockPersonality=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes

StandardOutput=journal
StandardError=journal
SyslogIdentifier=hart-llm

[Install]
WantedBy=hart.target
