# deploy/linux/systemd/hyve-vision.service
[Unit]
Description=HyveOS Vision Service (MiniCPM)
Documentation=https://github.com/hevolve-ai/hyve
After=hyve-backend.service
ConditionPathExists=/opt/hyve/models/minicpm/

[Service]
Type=simple
User=hyve
Group=hyve
WorkingDirectory=/opt/hyve
ExecStart=/opt/hyve/venv/bin/python integrations/vision/minicpm_server.py --model_dir /opt/hyve/models/minicpm --port ${HYVE_VISION_PORT:-9891} --device auto
EnvironmentFile=/etc/hyve/hyve.env
Restart=on-failure
RestartSec=15
TimeoutStartSec=120

# Vision service needs GPU access
SupplementaryGroups=video render
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/hyve /opt/hyve/models
PrivateTmp=yes

ProtectClock=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
SystemCallFilter=@system-service
MemoryDenyWriteExecute=no
LockPersonality=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes

StandardOutput=journal
StandardError=journal
SyslogIdentifier=hyve-vision

[Install]
WantedBy=hyve.target
