# deploy/linux/systemd/hyve-discovery.service
[Unit]
Description=HyveOS Peer Discovery (UDP Beacon)
Documentation=https://github.com/hevolve-ai/hyve
After=hyve-backend.service
BindsTo=hyve-backend.service

[Service]
Type=simple
User=hyve
Group=hyve
WorkingDirectory=/opt/hyve
ExecStart=/opt/hyve/venv/bin/python -c "from integrations.social.peer_discovery import AutoDiscovery; d = AutoDiscovery(); d.start(); import time; time.sleep(999999)"
EnvironmentFile=/etc/hyve/hyve.env
Restart=on-failure
RestartSec=10

# Needs UDP broadcast capability
AmbientCapabilities=CAP_NET_BROADCAST
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/hyve
PrivateTmp=yes

ProtectClock=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
SystemCallFilter=@system-service
MemoryDenyWriteExecute=no
LockPersonality=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes

StandardOutput=journal
StandardError=journal
SyslogIdentifier=hyve-discovery

[Install]
WantedBy=hyve.target
