# deploy/linux/systemd/hyve-agent-daemon.service
[Unit]
Description=HyveOS Agent Daemon (Goal Engine)
Documentation=https://github.com/hevolve-ai/hyve
After=hyve-backend.service
BindsTo=hyve-backend.service

[Service]
Type=simple
User=hyve
Group=hyve
WorkingDirectory=/opt/hyve
ExecStart=/opt/hyve/venv/bin/python -c "from integrations.agent_engine.agent_daemon import AgentDaemon; d = AgentDaemon(); d.run_forever()"
EnvironmentFile=/etc/hyve/hyve.env
Restart=on-failure
RestartSec=30
TimeoutStartSec=30

NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/hyve /var/log/hyve /opt/hyve/agent_data
PrivateTmp=yes

StandardOutput=journal
StandardError=journal
SyslogIdentifier=hyve-agent-daemon

[Install]
WantedBy=hyve.target
