# deploy/distro/update/hart-update.timer
[Unit]
Description=HART OS Daily Update Check

[Timer]
OnCalendar=*-*-* 03:00:00
RandomizedDelaySec=3600
Persistent=true

[Install]
WantedBy=timers.target
